# .config/battery_notify.service
[Unit]
Description=Battery Low Notification

[Service]
Type=oneshot
ExecStart=/home/YOURUSERNAME/.local/bin/battery_notify.sh
